# src/usr/lib/udev/rules.d/61-gdm.rules
# DON'T disable Wayland when using the proprietary nvidia driver (line below should be commented out)
# DRIVER=="nvidia", RUN+="/usr/libexec/gdm-disable-wayland"
